# .dev/edge-native-base/container_config/s6-overlay/s6-rc.d/oauth2_proxy/run
#!/command/with-contenv sh
exec 2>&1
set -e
export LC_ALL=en_US.utf-8
export LANG=en_US.utf-8
cd /opt
exec edm -r /opt/_edm_oauth2_proxy run -- python -m oauth2_proxy.wsgi