# .dev/edge-native-base/container_config/s6-overlay/s6-rc.d/location_config/up
with-contenv
exec /bin/bash -c "if [ -z $JUPYTERHUB_SERVICE_PREFIX]; then export JUPYTERHUB_SERVICE_PREFIX='/'; fi && envsubst '${JUPYTERHUB_SERVICE_PREFIX}' < /opt/nginx/app.location.conf.template > /opt/nginx/app.location.conf"