# .dev/edge-native-base/container_config/s6-overlay/s6-rc.d/application/run
#!/command/with-contenv sh

# Main application entrypoint

exec 2>&1
cd /home/app
exec ./startup-script.sh
