# streamlit-app-example/s6/etc/s6-overlay/s6-rc.d/nginx/run
#!/bin/sh
exec 2>&1
echo "Starting nginx"
set -e
exec nginx