# .dev/edge-native-base/container_config/s6-overlay/s6-rc.d/nginx/run
#!/command/with-contenv sh
exec 2>&1
echo "Starting nginx"
set -e
exec nginx -c /opt/nginx/nginx.conf