# edge-oauth2-app/container_config/etc/s6-overlay/s6-rc.d/application/run
#!/command/with-contenv sh

# Main application entrypoint

exec 2>&1
cd /opt
exec ./startup-script.sh
